# examples/launchd.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.todoist-to-obsidian.run</string>
  <key>ProgramArguments</key>
  <array>
    <string>/path/to/python</string>
    <string>-m</string>
    <string>todoist_obsidian_bridge</string>
    <string>run</string>
    <string>--config</string>
    <string>/path/to/config.toml</string>
  </array>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
